# geo-proxy-relay.service
[Unit]
Description=GEO Proxy Relay (FastAPI)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
WorkingDirectory=/opt/geo-proxy-relay
EnvironmentFile=/opt/geo-proxy-relay/.env
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
