# relay-proxy.service
[Unit]
Description=Relay Proxy (FastAPI)
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/opt/geo-proxy-relay
EnvironmentFile=/opt/geo-proxy-relay/.env
ExecStart=/opt/geo-proxy-relay/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target